# deploy/systemd/fusefs.service
[Unit]
Description=VirtualFS FUSE read-only mount
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/virtualfs/fusefs.env
ExecStart=/usr/local/bin/fusefs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
